# retailco-app-repo/deploy/retailco-app.service
[Unit]
Description=RetailCo order-management demo app
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/azureuser/azure-compute/retailco-app-repo
ExecStart=/usr/bin/node /home/azureuser/azure-compute/retailco-app-repo/server.js
Restart=always
RestartSec=5
Environment=PORT=8080

[Install]
WantedBy=multi-user.target
